# .config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper using jjk images

[Service]
Type=oneshot
ExecStart=/usr/bin/zsh $HOME/.config/scripts/jjk_wallpaper.zsh

